# etc/cron.d/cron-create-collage
# create e-helvetica collage
21 21 * * 6 cronuser /bin/update-collage.sh > /data/logs/cron/update-collage.log 2>&1
