# etc/cron.d/cron-sync-collage
# sync e-helvetica collage
01 01 * * 0 cronuser /bin/sync-collage.sh > /target/data/logs/cron/sync-collage.log 2>&1
